# ops/systemd/big5quest.service
[Unit]
Description=Big5-Quest Rails API (Puma)
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/opt/Big5-Quest

# 本番環境とBundler
Environment=RAILS_ENV=production
Environment=BUNDLE_WITHOUT=development:test
Environment=BUNDLE_PATH=/opt/Big5-Quest/vendor/bundle

# ★ ここからSECRET_KEY_BASE等を読む
EnvironmentFile=/etc/big5quest.env

# 初回起動で pid が無くて落ちないように
ExecStartPre=/usr/bin/env bash -lc 'mkdir -p tmp/pids'

# Puma起動
ExecStart=/usr/bin/env bash -lc 'bundle exec puma -b tcp://127.0.0.1:3000 -e ${RAILS_ENV} --pidfile tmp/pids/puma.pid'

Restart=always
RestartSec=3
TimeoutStopSec=30
PIDFile=/opt/Big5-Quest/tmp/pids/puma.pid

[Install]
WantedBy=multi-user.target