# ops/systemd/big5quest-weekly-rollover.service
[Unit]
Description=Big5-Quest weekly rollover
After=network.target big5quest.service

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/opt/Big5-Quest
Environment=RAILS_ENV=production
EnvironmentFile=/etc/big5quest.env
ExecStart=/usr/bin/env bundle exec rake weekly:rollover
